# dnprogs/debian/dnet-common.conffiles
/etc/decnet.proxy
/etc/default/decnet
/etc/dnetd.conf
/etc/init.d/decnet
/etc/modutils/decnet
